# config/systemd/user/waybar.service.d/override.conf
[Service]
Environment="GTK_DEBUG=interactive"
